# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:smbd]
command=/usr/sbin/smbd --daemon --foreground --log-stdout
redirect_stderr=true
